# systemd/bacalhau.service
[Unit]
Description=CoopHive Bacalhau
After=network-online.target
Wants=network-online.target systemd-networkd-wait-online.service

[Service]
Environment="LOG_TYPE=json"
Environment="LOG_LEVEL=debug"
Environment="HOME=/app/generic-dcn"
Environment="BACALHAU_SERVE_IPFS_PATH=/app/generic-dcn/ipfs"
Restart=always
RestartSec=5s
ExecStart=/usr/bin/bacalhau serve --node-type compute,requester --peer none --private-internal-ipfs=false --job-selection-accept-networked

[Install]
WantedBy=multi-user.target
